# crontab
*/10 * * * *  /home/pi/metar-map/start.sh
0 3 * * sun   /home/pi/metar-map/upgrade.sh
