# nodejs/crontab
*/30 * * * *    sudo /home/pi/metar-map/start.sh 2>&1 | /usr/bin/logger -t metarmap
0 3 * * sun     sudo /home/pi/metar-map/upgrade.sh 2>&1 | /usr/bin/logger -t metarmap
